# infra/systemd/haraka-inbound.service
[Unit]
Description=Inboxi Haraka inbound MTA (catch-all SMTP on port 25)
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/inboxi/apps/mta-inbound
# Source the repo-root .env (bash strips quotes correctly — systemd
# EnvironmentFile does not) then exec Haraka. Runs as root so it can bind
# port 25 and read .env without extra capability/permission setup.
ExecStart=/bin/bash -c 'set -a; . /opt/inboxi/.env; set +a; exec /opt/inboxi/apps/mta-inbound/node_modules/.bin/haraka -c /opt/inboxi/apps/mta-inbound'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
