# infra/systemd/haraka-outbound.service
[Unit]
Description=Inboxi Haraka outbound MTA (self-host delivery driver)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/inboxi/apps/mta-outbound
ExecStart=/bin/bash -c 'set -a; . /opt/inboxi/.env; set +a; exec /opt/inboxi/apps/mta-outbound/node_modules/.bin/haraka -c /opt/inboxi/apps/mta-outbound'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
