# infra/systemd/haraka-outbound.service
[Unit]
Description=Inboxi Haraka outbound MTA (self-host delivery driver)
After=network.target

[Service]
Type=simple
User=inboxi
WorkingDirectory=/opt/inboxi/apps/mta-outbound
# Run the Haraka bin shim directly (it is a shell wrapper, not a JS file).
ExecStart=/opt/inboxi/apps/mta-outbound/node_modules/.bin/haraka -c /opt/inboxi/apps/mta-outbound
Restart=always
RestartSec=5
EnvironmentFile=/opt/inboxi/.env

[Install]
WantedBy=multi-user.target
